# configs/system/hourly/system.hourly.timer
[Unit]
Description=Hourly Tasks

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target