# configs/system/hourly/system.hourly.service
[Unit]
Description=Hourly Script

[Service]
Type=simple
ExecStart=/var/system/scripts/hourly.sh
Slice=flexSystem.slice