# configs/system/logging/system.logDelete.service
[Unit]
Description=Flex Delete Logs

[Service]
Type=simple
ExecStart=/var/system/scripts/logging/logDelete.sh
WorkingDirectory=/tmp/log

