# configs/system/logging/system.logDelete.service
[Unit]
Description=Flex Delete Logs

[Service]
Type=simple
ExecStart=/home/flex/scripts/flexLogDelete.sh
WorkingDirectory=/tmp/log

