# configs/gui/systemd/gui.slice
[Unit]
Description=Gui Slice
DefaultDependencies=no

[Slice]